# deploy/systemd/web_demo.service.template
[Unit]
Description=Lab Safe Assistant Web Demo
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_DIR__
EnvironmentFile=__ENV_FILE__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PYTHON_BIN__ -m uvicorn web_demo.app:app --host 0.0.0.0 --port __DEMO_PORT__
Restart=always
RestartSec=3
TimeoutStopSec=15

[Install]
WantedBy=default.target
